# topology-cacher.cron
SHELL=/bin/bash
0,15,30,45 * * * * root sleep $(( $RANDOM \% 180 )) && /usr/bin/python3 -m topology_cacher
