# docker/supervisor-apache.conf
[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
autorestart=true
