# docker/supervisor-varnish.conf
[program:varnish]
# Varnish listening on port 8000
command=/usr/sbin/varnishd -F -a :8000 -f /etc/varnish/topology.vcl
autorestart=true
