# docker/supervisor-apache.conf
[program:httpd]
command=/bin/bash -c "source /etc/sysconfig/httpd && exec /usr/sbin/httpd $OPTIONS -DFOREGROUND"
autorestart=true
